# gentoo/etc/init.d/ycast/ycast
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

#extra_started_commands="reload"
command="/usr/local/bin/ycast.py"
command_args="${YCAST_ARGS}"
description="Fake vtuner Radio Directory Service"

depend() {
	need net
	use localmount
}

